# monitor/launchd/com.srijan.market-pulse.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- INSTALLED 2026-07-06 (Srijan approved after reviewing the first pulse draft).
     PLACEMENT justification (per ~/automation-migration/PLACEMENT.md): laptop LaunchAgent
     because the job writes/reads the LOCAL aksh market-pulse SQLite (data/market_pulse_history.db)
     and its drafts feed the laptop-bound Chrome-CDP X posting lane; no cloud copy of the aksh
     data layer exists. ONE job chains producer + draft + video spec to avoid agent sprawl (the
     61-agent lesson) AND the draft-before-converter race that killed the wrap lane on 2026-07-06:
     the video converter runs strictly after its draft exists. Replaces com.srijan.market-video
     (16:10, retired 2026-07-06 — its wrap-draft writer was never scheduled and died 07-03).
     16:20 IST: after NSE close. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.srijan.market-pulse</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd "$HOME/aksh_backtesting_trading" &amp;&amp; /opt/homebrew/bin/python3 -m scripts.market_pulse_daily_snapshot &amp;&amp; cd "$HOME/loop engineering/content-engine" &amp;&amp; /opt/homebrew/bin/python3 monitor/market_pulse_post.py &amp;&amp; cd "$HOME/content creation" &amp;&amp; /opt/homebrew/bin/python3 publish/market_video.py --lane pulse</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>16</integer>
    <key>Minute</key>
    <integer>20</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/srijan/autonomy/logs/market-pulse.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/srijan/autonomy/logs/market-pulse.err</string>
</dict>
</plist>
