# monitor/launchd/com.srijan.backtest-education.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- INSTALLED 2026-07-11 (first draft reviewed+posted 2026-07-06; Srijan's go 2026-07-11).
     Dispatches via x/cron/run.sh so the weekly draft is also POSTED (generator alone left
     drafts unposted). Off-market-hours (Sun) so it never competes with the daily finance
     trio or the market-pulse job; post window pinned in LANE_WINDOWS["backtest-education"]. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.srijan.backtest-education</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/srijan/loop engineering/content-engine/x/cron/run.sh</string>
    <string>backtest-education</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/srijan/autonomy/logs/backtest-education.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/srijan/autonomy/logs/backtest-education.err</string>
</dict>
</plist>
